# gadget.service
[Unit]
Description=My USB gadget service

[Service]
ExecStart=/usr/bin/gadget.sh

[Install]
WantedBy=multi-user.target

